# init.d/pzstd.service
[Service]
ExecStart=/usr/local/bin/pzstd
Restart=always

[Install]
WantedBy=default.target
